# ansible/roles/celery/templates/celery-cms.target.j2
[Unit]
Description=Celery CMS Workers Target
After=network.target

[Install]
WantedBy=multi-user.target
