# ansible/roles/celery/templates/celery-worker@.service.j2
[Unit]
Description=Celery worker for CMS - queue {{ queue_name }}
After=network.target redis-server.service
PartOf=celery-cms.target

[Service]
Type=simple
User={{ deploy_user }}
Group={{ deploy_group }}
WorkingDirectory={{ app_path }}
Environment="PATH={{ app_path }}/venv/bin"
Environment="DJANGO_SETTINGS_MODULE={{ django_settings_module }}"
Environment="CELERY_BROKER_URL={{ celery_broker_url }}"
Environment="CELERY_RESULT_BACKEND={{ celery_result_backend }}"

ExecStart={{ app_path }}/venv/bin/celery \
  -A cms.celery.app worker \
  -Q {{ queue_name }} \
  --loglevel={{ queue_config.loglevel }} \
  --concurrency={{ queue_config.concurrency }} \
  --hostname=worker-{{ queue_name }}@%h

ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
